# etc/periodic/daily/update-ntp
#!/bin/sh
ntpd -d -q -n -p pool.ntp.org
